# anisystemd.service
[Unit]
Description=AniDev Discord Bot Service (Systemd with Auto-Restart)
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/anisystemd
Restart=always
RestartSec=10
WatchdogSec=60

# Переменные окружения
Environment="DISCORD_TOKEN=your_discord_token_here"
EnvironmentFile=-/etc/anicore/env

# Безопасность
User=anicore
Group=anicore
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/anicore

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anisystemd

# Ограничения ресурсов
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

